# docs/systemd/vroom.service
[Unit]
Description=VROOM daemon
After=mariadb.service

[Service]
Type=simple
ExecStart=/opt/vroom/scripts/start_server.pl
PrivateTmp=yes
User=vroom
Group=vroom
SyslogIdentifier=vroom
InaccessibleDirectories=/home
ReadOnlyDirectories=/usr
ReadOnlyDirectories=/etc
Restart=on-failure

[Install]
WantedBy=multi-user.target
Wants=mariadb.service httpd.service

